# ii-bot/supervisor.conf
{% for server_address, server in bot.get("servers", {}).items() %}
[program:{{ bot_name }}_{{ server_address|replace(".", "_") }}_reconnection_loop]
command=bash reconnection_{{ server_address|replace(".", "_") }}.sh
directory=/var/ii/{{ bot_name }}/
user=ii
redirect_stderr=True
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true

{% for rss in server.get("rss", []) %}
[program:{{ bot_name }}_{{ server_address|replace(".", "_") }}_rss_{{ rss_name }}]
command=bash rss_{{ rss_name }}_{{ server_address|replace(".", "_") }}.sh
directory=/var/ii/{{ bot_name }}/
user=ii
redirect_stderr=True
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true

{% endfor %}

{% for search in server.get("twitter_search", []) %}
[program:{{ bot_name }}_{{ server_address|replace(".", "_") }}_twitter_search_{{ search.name }}]
command=bash twitter_search_{{ search.name }}_{{ server_address|replace(".", "_") }}.sh
directory=/var/ii/{{ bot_name }}/
user=ii
redirect_stderr=True
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true

{% endfor %}
{% endfor %}
