# ii-bot/supervisor.conf
[program:{{ bot.name }}_reconnection_loop]
command=./reconnection.sh
directory=/var/ii/{{ bot.name }}/
user=ii
redirect_stderr=True
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true
